# ddns.service
[Unit]
Description=Ddns service

[Service]
Type=simple
StateDirectory=simpleddns
ExecStart=/usr/sbin/ddns-notify
TimeoutStartSec=65s
TimeoutStopSec=30s


[Install]
WantedBy=multi-user.target