# ddns.service
[Unit]
Description=Ddns service

[Service]
Type=simple
StateDirectory=simpleddns
ExecStart=/usr/bin/simpleddns
TimeoutStartSec=65s
TimeoutStopSec=30


[Install]
WantedBy=multi-user.target